# conf/aethercast.conf.in
description "Display cast service"

start on started dbus
stop on stopping dbus

respawn

exec /usr/sbin/aethercast
