# conf/miracast-service.conf.in
description "Miracast Management Service"

start on started dbus
stop on stopping dbus

respawn

exec /usr/sbin/miracast-service
